# root/etc/services.d/wireguard/run
#!/usr/bin/with-contenv bash

_term() {
  echo "Caught SIGTERM signal!"
  wg-quick down ${DEV_NAME}
}

trap _term SIGTERM

wg-quick up ${DEV_NAME}

sleep infinity &

wait
