# root/etc/s6-overlay/s6-rc.d/svc-wireguard/run
#!/usr/bin/with-contenv bash

wg-quick up ${IF_NAME}
